# rpm/redhat/9.5/pgpool-II/EL-5/pgpool-II-94.service
[Unit]
Description=PGPool-II Middleware Between PostgreSQL Servers And PostgreSQL Database Clients
After=syslog.target network.target

[Service]
EnvironmentFile=-/etc/sysconfig/pgpool-II-94
ExecStart=/usr/pgpool-9.4/bin/pgpool -f /etc/pgpool-II-94/pgpool.conf $OPTS 
ExecStop=/usr/pgpool-9.4/bin/pgpool -f /etc/pgpool-II-94/pgpool.conf -m fast stop

[Install]
WantedBy=multi-user.target
